# crontab
30 2 * * * ./Monitor/cron_script.sh > /root/Monitor/crontab.log